# claude-pm-doc-sync.service
[Unit]
Description=Claude Multiagent PM Automated Documentation Synchronization Service
After=network.target

[Service]
Type=simple
User=masa
WorkingDirectory=/Users/masa/Projects/claude-multiagent-pm
ExecStart=/usr/bin/python3 /Users/masa/Projects/claude-multiagent-pm/scripts/automated_doc_sync.py --daemon
Restart=always
RestartSec=10
StandardOutput=append:/Users/masa/Projects/claude-multiagent-pm/logs/automated_doc_sync.log
StandardError=append:/Users/masa/Projects/claude-multiagent-pm/logs/automated_doc_sync_error.log

[Install]
WantedBy=multi-user.target
