# claude-pm-health-monitor.service
[Unit]
Description=Claude Multiagent PM Framework Health Monitor
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=masa
WorkingDirectory=/Users/masa/Projects/claude-multiagent-pm
ExecStart=/usr/local/bin/node scripts/automated-health-monitor.js monitor --interval=5 --threshold=60
StandardOutput=append:/Users/masa/Projects/claude-multiagent-pm/logs/health-monitor-systemd.log
StandardError=append:/Users/masa/Projects/claude-multiagent-pm/logs/health-monitor-systemd-error.log
SyslogIdentifier=claude-multiagent-pm-health-monitor
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target