# enhanced-claude-pm-doc-sync.cron
# Enhanced Claude PM Documentation Synchronization
# Run every 5 minutes with enhanced validation
*/5 * * * * python3 /Users/masa/Projects/claude-multiagent-pm/scripts/enhanced_automated_doc_sync.py --once >> /Users/masa/Projects/claude-multiagent-pm/logs/cron_enhanced_doc_sync.log 2>&1
